# deploy/systemd/aurum-cex-bot-loop.service
[Unit]
Description=Aurum CEX resident mechanical paper bot loop
Wants=network-online.target aurum-cex-recorder.service
After=network-online.target aurum-cex-recorder.service

[Service]
Type=simple
User=aurum
Group=aurum
WorkingDirectory=/opt/aurum/app
EnvironmentFile=/opt/aurum/.env
Environment=AURUM_ENV_FILE=/opt/aurum/.env
Environment=AURUM_CEX_DATA_DIR=/opt/aurum/data/cex_arena
Environment=AURUM_DATA_DIR=/opt/aurum/data/paper_duel
Environment=AURUM_PUBLIC_DASHBOARD_DIR=/opt/aurum/public/dashboard
Environment=AURUM_CEX_STRATEGIES=momentum,mean_reversion,volatility_breakout,funding_filter,grid_shadow
Environment=AURUM_CEX_NOTIONAL_USDT=50
Environment=AURUM_CEX_MAX_ORDERS=1
Environment=AURUM_CEX_MAX_STALE_SECONDS=420
Environment=AURUM_CEX_BOT_INTERVAL_SEC=30
Environment=AURUM_CEX_ENABLE_LIVE_TRADING=false
Environment=AURUM_RUN_DIR=/opt/aurum/run
Environment=AURUM_LOG_DIR=/opt/aurum/logs
ExecStart=/opt/aurum/app/scripts/run_cex_bot_loop.sh
Restart=always
RestartSec=10
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
ReadWritePaths=/opt/aurum/data /opt/aurum/logs /opt/aurum/run /opt/aurum/public/dashboard

[Install]
WantedBy=multi-user.target
