# deploy/systemd/aurum-strategy-review.timer
[Unit]
Description=Run Aurum slow strategy review about every 30 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
Persistent=true
RandomizedDelaySec=180
Unit=aurum-strategy-review.service

[Install]
WantedBy=timers.target
