# deploy/systemd/aurum-strategy-review.service
[Unit]
Description=Aurum slow strategy review and gated rule prompt update
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=aurum
Group=aurum
WorkingDirectory=/opt/aurum/app
EnvironmentFile=/opt/aurum/.env
Environment=AURUM_ENV_FILE=/opt/aurum/.env
Environment=AURUM_DATA_DIR=/opt/aurum/data/paper_duel
Environment=AURUM_RUN_DIR=/opt/aurum/run
Environment=AURUM_LOG_DIR=/opt/aurum/logs
Environment=AURUM_PUBLIC_DASHBOARD_DIR=/opt/aurum/public/dashboard
Environment=AURUM_REVIEW_CADENCE_SECONDS=1800
ExecStart=/opt/aurum/app/scripts/run_strategy_review.sh
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
ReadWritePaths=/opt/aurum/data /opt/aurum/logs /opt/aurum/run /opt/aurum/public/dashboard
