# deploy/systemd/aurum-cex-recorder.service
[Unit]
Description=Aurum CEX public-data recorder (paper-only)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=aurum
Group=aurum
WorkingDirectory=/opt/aurum/app
Environment=APP_DIR=/opt/aurum/app
Environment=AURUM_CEX_DATA_DIR=/opt/aurum/data/cex_arena
Environment=AURUM_CEX_MARKETS=okx:BTC-USDT-SWAP,okx:ETH-USDT-SWAP
Environment=AURUM_CEX_RECORDER_INTERVAL_SEC=30
Environment=AURUM_CEX_ENABLE_LIVE_TRADING=false
Environment=AURUM_RUN_DIR=/opt/aurum/run
Environment=AURUM_LOG_DIR=/opt/aurum/logs
ExecStart=/opt/aurum/app/scripts/run_cex_recorder.sh
Restart=always
RestartSec=10
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aurum/data /opt/aurum/logs /opt/aurum/run
ReadOnlyPaths=/opt/aurum/app
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
