# deploy/systemd/aurum-bot-loop.service
[Unit]
Description=Aurum resident mechanical paper bot loop
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=aurum
Group=aurum
WorkingDirectory=/opt/aurum/app
EnvironmentFile=/opt/aurum/.env
Environment=AURUM_ENV_FILE=/opt/aurum/.env
Environment=AURUM_DATA_DIR=/opt/aurum/data/paper_duel
Environment=AURUM_RUN_DIR=/opt/aurum/run
Environment=AURUM_LOG_DIR=/opt/aurum/logs
Environment=AURUM_PUBLIC_DASHBOARD_DIR=/opt/aurum/public/dashboard
Environment=AURUM_BOT_MIN_INTERVAL_SEC=5
Environment=AURUM_BOT_DEFAULT_INTERVAL_SEC=15
ExecStart=/opt/aurum/app/scripts/run_bot_loop.sh
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
ReadWritePaths=/opt/aurum/data /opt/aurum/logs /opt/aurum/run /opt/aurum/public/dashboard

[Install]
WantedBy=multi-user.target
